# systemd/livepi-video-glitch.service
# Installed on the Pi as /etc/systemd/system/livepi-video-glitch.service.
# See docs/deploy.md for the full setup (X11 kiosk session, group
# permissions, disabling any conflicting auto-login desktop). This unit is
# designed but not yet empirically verified against real hardware -- Phase 6
# in docs/architecture.md.
#
# Edit User=/Group= below if not running as the default "pi" user.

[Unit]
Description=LivePi VideoGlitcher (kiosk)
After=sound.target network.target

[Service]
Type=simple
EnvironmentFile=/etc/livepi-video-glitch.env
ExecStart=/bin/sh -c 'cd "$APP_DIR" && exec startx "$APP_DIR/bin/livepi-video-glitch" -- -s off -dpms'
Restart=on-failure
RestartSec=2
User=pi
Group=pi
SupplementaryGroups=video render audio

[Install]
WantedBy=multi-user.target
